# rootfs/etc/service/node/run
#!/bin/sh -e

# pipe stderr to stdout and run node server
exec 2>&1
exec node /var/www/html/index.js